# hylafax+-srpm/hylafax+_hourly.cron
#!/bin/sh

if [ -f /etc/sysconfig/hylafax+ ]; then
  . /etc/sysconfig/hylafax+
fi

if [ "$FAXQCLEAN_RUN" != "no" ]; then
  /usr/sbin/faxqclean $FAXQCLEAN_OPTIONS
fi

exit 0
